# relay/deploy/serpmonn-relay.service
[Unit]
Description=Serpmonn P2P Relay/Bootstrap Node
Documentation=https://github.com/serpmonn/Serpmonn
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=serpmonn
Group=serpmonn

# UDP/443 (QUIC) без root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# Бинарник: /usr/local/bin/serpmonn-relay
# Снаружи только 443: WS на 127.0.0.1:4002 за nginx (/api/v1/sync),
# технические 4001 — тоже localhost.
ExecStart=/usr/local/bin/serpmonn-relay \
    -key /etc/serpmonn/relay-identity.key \
    -port 4001 \
    -wsport 4002 \
    -port443 443 \
    -extip 188.235.13.20 \
    -domain relay.serpmonn.ru

Restart=always
RestartSec=5s

LimitNOFILE=65536
LimitNPROC=4096

StandardOutput=journal
StandardError=journal
SyslogIdentifier=serpmonn-relay

[Install]
WantedBy=multi-user.target
